# crontab
*/5 * * * * /usr/bin/python clear_cache.py
