# cronjob/crontab
*/1 * * * * /clear_cache.sh >> /var/log/cron.log 2>&1
