# cronjob/crontab
0 0 * * * /clear_cache.sh >> /var/log/cron.log 2>&1
